# deploy/openclaw/mvm-tick-alert.service
[Unit]
Description=Mvm tick failure alert

[Service]
Type=oneshot
ExecStart=/usr/bin/logger -t mvm-alert "mvm-tick failed at $(date -u +%%Y-%%m-%%dT%%H:%%M:%%SZ)"
